# files/budgie/product/etc/xdg/autostart/plank.desktop
[Desktop Entry]
Name=Plank
Comment=Stupidly simple.
Exec=plank
Icon=plank
Terminal=false
Type=Application
NoDisplay=true
